# start.conf
description     "Start Cuckoo Stuff"
start on startup
task
exec /etc/
